# systemd/user/xinit.service
[Unit]
Description=Xinit

[Service]
ExecStart=/usr/bin/xinit -- :0 vt7

[Install]
WantedBy=wm.target
